# 411-cron/config/411.cron
* * * * * /usr/bin/php /var/www/411/bin/cron.php >>  > /proc/1/fd/1 2>/proc/1/fd/2
* * * * * /usr/bin/php /var/www/411/bin/worker.php >>  > /proc/1/fd/1 2>/proc/1/fd/2
